# systemd/surge_xt_wrapper.service
[Unit]
Description=Start Surge XT Wrapper
After=graphical.target

[Service]
User=root
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/pi/.Xauthority
ExecStart=/home/pi/microtonOS/surge_xt_wrapper.py
Restart=on-failure
RestartSec=2s
KillMode=control-group
TimeoutSec=infinity

[Install]
WantedBy=graphical.target
